# systemd/multi-casino-scraper.timer
[Unit]
Description=Run multi-casino scraper every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
