# systemd/multi-casino-scraper.service
[Unit]
Description=Multi-Casino Jackpot Aggregator
After=network.target postgresql.service

[Service]
Type=oneshot
User=rod
WorkingDirectory=/home/rod/home_ai_stack/jackpot-ingest
ExecStart=/home/rod/home_ai_stack/venv/bin/python3 scripts/aggregate_all.py
StandardOutput=journal
StandardError=journal

# Restart on failure
Restart=on-failure
RestartSec=120s
